# agenticbro-backend-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agenticbro.backend-server</string>
  
  <key>ProgramArguments</key>
  <array>
    <string>/Users/efinney/.nvm/versions/node/v25.6.1/bin/node</string>
    <string>/Users/efinney/.nvm/versions/node/v25.6.1/bin/npx</string>
    <string>tsx</string>
    <string>watch</string>
    <string>server/index.ts</string>
  </array>
  
  <key>WorkingDirectory</key>
  <string>/Users/efinney/.openclaw/workspace/aibro</string>
  
  <key>RunAtLoad</key>
  <true/>
  
  <key>KeepAlive</key>
  <true/>
  
  <key>StandardOutPath</key>
  <string>/Users/efinney/.openclaw/logs/agenticbro-backend-server.log</string>
  
  <key>StandardErrorPath</key>
  <string>/Users/efinney/.openclaw/logs/agenticbro-backend-server.err.log</string>
  
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/Users/efinney/.nvm/versions/node/v25.6.1/bin</string>
  </dict>
</dict>
</plist>